# my_crontab/drink.crontab
*/30 10-18 * * 1-5 ~/.config/my_crontab/drink_water.sh
* 15 * * 1-5 ~/.config/my_crontab/drink_tea.sh
* 19 * * 1-5 ~/.config/my_crontab/drink_beer.sh
