# my_crontab/drink_water.crontab
*/30 10-18 * * 1-5 ~/.config/my_crontab/drink_water.sh
